# .devcontainer/build/crontab/root
* * * * * cd /workspace/example-app && php artisan schedule:run >> /var/log/cron.log
* * * * * echo "Hello" >> /workspace/cron.txt
